# usr/local/jrootfs/etc/init.d/telnetd
#!/sbin/openrc-run

description="telnet server"

procname="telnetd"
command="/usr/sbin/telnetd"
command_args="-b 127.0.0.1:9923 -K -l/bin/sh"

depend() {
	use logger dns
	after loopback
}
